# nanoclaw.service
[Unit]
Description=NanoClaw - Personal Claude Assistant
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=%I
WorkingDirectory=/home/%I/nanoclaw
Environment="NODE_ENV=production"
Environment="HOME=/home/%I"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# IMPORTANT: Update the following paths according to your Node.js installation
# For NodeSource installation (default): /usr/bin/node and /usr/bin/npm
# For nvm installation: /home/%I/.nvm/versions/node/vXX.X.X/bin/node
ExecStart=/usr/bin/node /usr/bin/npm run start

# Restart policy
Restart=always
RestartSec=10

# Graceful shutdown
TimeoutStopSec=30
KillSignal=SIGTERM

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nanoclaw

[Install]
WantedBy=multi-user.target
